# docker/supervisord.conf
[supervisord]
user = root
nodaemon = true
debug = debug

[program:hhvm]
command = hhvm -m server -d hhvm.server.type=fastcgi -d hhvm.server.port=9000

[program:nginx]
command = nginx
